# ansible/roles/analytics/roles/zeppelin/files/etc/profile.d/zeppelin.sh
#!/usr/bin/env bash
export ZEPPELIN_HOME=/opt/zeppelin/zeppelin-latest

alias ll='ls -al'
alias zconf='cd $ZEPPELIN_HOME/conf; ll'
alias zlogs='cd $ZEPPELIN_HOME/logs; ll'