# ansible/roles/cloud/gcp/network/files/etc/network/if-up.d/nat-gateway
#!/bin/sh
# Copyright 2018, Dematic, Corp.

sudo sh -c "echo 1 > /proc/sys/net/ipv4/ip_forward"
# Masquerade.
sudo iptables -t nat -A POSTROUTING -o eth0 -j MASQUERADE

