# ansible/roles/common/files/etc/profile.d/dematiclabs.sh
#!/usr/bin/env bash

# Copyright 2018 Dematic, Corp.  # Licensed under the MIT Open Source License: https://opensource.org/licenses/MIT

export JAVA_HOME=/opt/jdk_latest

alias countfiles='find . -type f | wc -l'