# ansible/roles/analytics/roles/producer/files/etc/profile.d/dlabs-producer.sh
#!/usr/bin/env bash
export JAVA_HOME=/opt/jdk_latest

export PS1="[\u@\H \w]"

